# nginx-consul/s6/consul-template/run
#!/bin/sh

sleep 3
exec consul-template -config=/consul-template/config
