# nginx-consul/s6/consul-template/run
#!/bin/sh

# Sleep 3sec to allow nignx to start before changing the config again

sleep 3
exec consul-template -config=/consul-template/config
