# couchdb-proxy-legacy/s6/nginx/run
#!/bin/sh

rm -f /etc/nginx/conf.d/default.conf

timeout 20 ./wait4conf.sh || exit 1

exec nginx -g "daemon off;"
