# nginx-consul/s6/nginx/run
#!/bin/sh

# Run consul-template once to populate nginx config
consul-template -config=/consul-template/config -once

# Pass execution to nginx
exec nginx
