# sensu-aio/supervisord.ini
[supervisord]
nodaemon=true

[unix_http_server]
file=/var/run//supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run//supervisor.sock

[program:rabbitmq]
priority=10
environment=HOME="/var/lib/rabbitmq"
directory=/tmp
command=/usr/sbin/rabbitmq-server
user=rabbitmq
autostart=true
autorestart=true
stopsignal=QUIT

[program:redis]
priority=10
directory=/tmp
command=redis-server
user=redis
autostart=true
autorestart=true

[program:sensu-server]
priority=20
directory=/tmp
environment=PATH="/opt/sensu/embedded/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/etc/sensu/plugins:/etc/sensu/handlers",GEM_PATH="/opt/sensu/embedded/lib/ruby/gems/2.0.0",LANG="en_AU.utf8"
command=/opt/sensu/bin/sensu-server -c /etc/sensu/config.json -d /etc/sensu -e /etc/sensu/extensions -v -l /var/log/sensu/server.log
user=sensu
startsecs=5
autostart=true
autorestart=true

[program:sensu-api]
priority=30
directory=/tmp
environment=PATH="/opt/sensu/embedded/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/etc/sensu/plugins:/etc/sensu/handlers",GEM_PATH="/opt/sensu/embedded/lib/ruby/gems/2.0.0",LANG="en_AU.utf8"
command=/opt/sensu/bin/sensu-api -c /etc/sensu/config.json -d /etc/sensu -e /etc/sensu/extensions -v -l /var/log/sensu/api.log
user=sensu
startsecs=5
autostart=true
autorestart=true

[program:sensu-client]
priority=40
directory=/tmp
environment=PATH="/opt/sensu/embedded/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/etc/sensu/plugins:/etc/sensu/handlers",GEM_PATH="/opt/sensu/embedded/lib/ruby/gems/2.0.0",LANG="en_AU.utf8"
command=/opt/sensu/bin/sensu-client -c /etc/sensu/config.json -d /etc/sensu -e /etc/sensu/extensions -v -l /var/log/sensu/client.log
user=sensu
autostart=true
autorestart=true

[program:uchiwa]
priority=40
directory=/opt/uchiwa/usr/src/uchiwa
command=/opt/uchiwa/embedded/bin/node app.js -c /etc/sensu/uchiwa.json
user=uchiwa
startsecs=5
autostart=true
autorestart=true

